# systemd/ansible-pull.service
[Unit]
Description=Apply desired-state configuration via ansible-pull
After=network-online.target
Wants=network-online.target
# Don't fight with package managers running in parallel.
ConditionACPower=true

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/ansible-pull
ExecStartPre=-/usr/bin/curl -fsSL ${REPO_URL}/raw/${REPO_BRANCH}/requirements.yml -o /tmp/requirements.yml
ExecStartPre=-/usr/bin/ansible-galaxy collection install -r /tmp/requirements.yml
ExecStart=/usr/bin/ansible-pull -U ${REPO_URL} -C ${REPO_BRANCH} -i localhost, local.yml

# Cap log noise on tiny boxes.
StandardOutput=journal
StandardError=journal
# Don't let a hung run block the next tick forever.
TimeoutStartSec=20m

[Install]
WantedBy=multi-user.target
